# roles/borg_backup/templates/backup.service.j2
[Unit]
Description=Backup {{ borg_backup_service_name }}

[Service]
User={{ borg_backup_user }}
Group={{ borg_backup_group }}
Type=simple
ExecStart=/usr/local/bin/borg_backup
