# roles/borg_backup/templates/backup.timer.j2
[Unit]
Description=Timer for {{ borg_backup_service_name }}

[Timer]
OnCalendar={{ borg_backup_calendar }}
Persistent=false

[Install]
WantedBy=timers.target
