# roles/backup_puller/templates/backup.timer.j2
[Unit]
Description=Backup {{ item.name }} timer

[Timer]
OnCalendar={{ item.calendar }}
Persistent=true

[Install]
WantedBy=timers.target
