# roles/backup_puller/templates/backups.target.j2
[Unit]
Description=Backups
After=network-online.target multi-user.target
AllowIsolate=no
DefaultDependencies=no

[Install]
WantedBy=multi-user.target
