# templates/rw-watchdog.timer.j2
[Unit]
Description=Reboot host when root filesystem becomes read-only

[Timer]
OnCalendar=minutely
Persistent=false

[Install]
WantedBy=timers.target
